# apps/gateway/deploy/bzz-ask-tls.service
[Unit]
Description=Swarm gateway on-demand TLS gatekeeper (caddy ask endpoint)
After=network.target

[Service]
Type=simple
# loopback-only Node responder; authorizes cert issuance for browserbzz.link
ExecStart=/usr/bin/node /var/www/bzz-gateway/deploy/ask-tls.js
Environment=GW_ZONE=browserbzz.link
Environment=GW_ASK_PORT=9111
Restart=on-failure
RestartSec=2
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
